# etc/systemd/user/xmonad.service
[Unit]
Description=restart xmonad on failure
PartOf=graphical-session.target

[Install]
WantedBy=xsession.target

[Service]
ExecStart=/home/hejohns/.cabal/bin/xmonad --replace
Type=simple
Restart=on-failure
